# deploy/f8-screener.service
[Unit]
Description=F8 F13 SEC Filing Screener
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=root
WorkingDirectory=/opt/F8_F13Screener
Environment=PYTHONUNBUFFERED=1
Environment=PYTHONPATH=/opt/F8_F13Screener
ExecStart=/opt/F8_F13Screener/venv/bin/python -m src.cli.main
Restart=always
RestartSec=15
StandardOutput=journal
StandardError=journal
# Give it a moment to finish in-flight Telegram sends before killing
TimeoutStopSec=30

[Install]
WantedBy=multi-user.target
